# fs/etc/services.d/ingress-controller/run
#!/usr/bin/with-contenv sh

exec /haproxy-ingress-controller --with-s6-overlay ${EXTRA_OPTIONS}
