# fs/etc/services.d/ingress-controller/run
#!/usr/bin/with-contenv sh

MEMLIMIT=$(free -m | awk '/Mem:/ {print int($2 / 3)}')
export GOMEMLIMIT="${MEMLIMIT}MiB"

exec /haproxy-ingress-controller --with-s6-overlay ${EXTRA_OPTIONS}
