# fs/etc/s6-overlay/s6-rc.d/ingress-controller/finish
#!/command/with-contenv sh

echo "Ingress Controller exited with code ${1}, taking down the S6 supervision tree"
exec /run/s6/basedir/bin/halt